# deploy/supervisord.conf
; dadaozhijian_cognitive — supervisord 配置（容器/裸机通用的进程监管备选）
; 适用：不想引入 systemd 的容器（PID 1 用 supervisord）或开发/边缘盒子。
;
; 用法：
;   pip install supervisor
;   supervisord -c deploy/supervisord.conf
;   supervisorctl -c deploy/supervisord.conf status

[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:dadaozhijian]
command=/usr/local/bin/dadaozhijian
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
environment=CONFIG_PATH="/etc/dadaozhijian/defaults.toml"
stdout_logfile=/var/log/dadaozhijian.out.log
stderr_logfile=/var/log/dadaozhijian.err.log
; 内存超限由 supervisor 的 memmon 事件（需 eventlistener）或外部监控处理

; 可选：指标导出定时任务（textfile 模式，免 HTTP）
; [program:dadaozhijian_metrics]
; command=/bin/sh -c 'while true; do mojo run -I . -I core ops/metrics_exporter.mojo /var/lib/node_exporter/dadaozhijian.prom; sleep 15; done'
; autostart=true
; autorestart=true
